# thirdparty/qt/automotivedemo/scripts/neptune.service
[Unit]
Description=Neptune
After=systemd-user-sessions.service

[Service]
ExecStart=/opt/automotivedemo/start.sh
